# flask-gunicorn-nginx/templates/default/gunicorn.conf.erb
[supervisord]
logfile=<%= File.join(helper.deploy_to, "log/supervisord.log") %>

[program:gunicorn]
command=<%= helper.gunicorn_start_path %>
directory=<%= @app_dir %>
autostart=true
autorestart=unexpected
stdout_logfile=<%= File.join(helper.deploy_to, "log/gunicorn.log") %>
redirect_stderr=true
