# flask-wsgi-nginx/templates/default/wsgi.conf.erb
[supervisord]
logfile=<%= File.join(helper.deploy_to, "log/supervisord.log") %>

[program:wsgi]
command=<%= helper.wsgi_start_path %>
directory=<%= @app_dir %>
autostart=true
autorestart=unexpected
stdout_logfile=<%= File.join(helper.deploy_to, "log/wsgi.log") %>
redirect_stderr=true

<% @memmon = node["flask-wsgi-nginx"]["supervisor"]["memmon"] %>

<% if @memmon["enabled"] %>
[eventlistener:memmon]
command=memmon --program wsgi=<%= @memmon["max_rss"] %>MB --cumulative --email=<%= @memmon["email"] %>
events=TICK_60
<% end %>
